# dotfiles/.config/systemd/user/clean-screenshots.service
[Unit]
Description=Clean up old screenshots

[Service]
Type=oneshot
ExecStart=/usr/bin/find %h/Screenshots -type f -mtime +7 -delete
